# shell/network-startup.service
[Unit]
Description=Network Startup Configuration
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/bin/sleep 10
ExecStart=/usr/local/bin/network-startup.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target